# service-scripts/start-honeyd.timer
[Unit]
Description=Honeyd file swap .

[Timer]
OnBootSec=9min
Unit=start-honeyd.service

[Install]
WantedBy=multi-user.target

